# etc/supervisor/conf.d/cron.conf
[program:cron]
environment=WTL_PRODUCTION=%(ENV_WTL_PRODUCTION)s
command=/usr/sbin/cron -f
