# etc/supervisor/conf.d/nginx.conf
[program:nginx]
environment=WTL_PRODUCTION=%(ENV_WTL_PRODUCTION)s
command=/usr/sbin/nginx -c /etc/nginx/nginx.conf
