# etc/supervisor/conf.d/php-fpm.conf
[program:php]
environment=WTL_PRODUCTION=%(ENV_WTL_PRODUCTION)s
user=www-data
group=www-data
command=/usr/sbin/php-fpm7.0 --nodaemonize --fpm-config /etc/php/7.0/fpm/php-fpm.conf
