# etc/cron.d/mediawiki-log-clean
0 * * * * root find /var/log/mediawiki/ -type f -mtime +7  -delete
